# deploy/collab/neko-collab.service
[Unit]
Description=Zhongheng browser collaboration service
Requires=docker.service
Wants=network-online.target
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=lamarck
Group=lamarck
SupplementaryGroups=docker
WorkingDirectory=/home/lamarck/services/neko-collab
ExecStartPre=/usr/bin/docker compose config --quiet
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecReload=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose stop
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
